# cron.d/boot-script
@reboot logbox /bin/echo "startup script works" > /tmp/message
