# cron.d/boot-script
@reboot logbox /bin/bash /home/logbox/os/boot
