# cron.d/report-ip
* * * * * logbox /bin/bash /home/logbox/os/report-ip.bash
